# server/systemd/codexwall.service
[Unit]
Description=Codex Limits JSON server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=codexwall
Group=codexwall
WorkingDirectory=/opt/codexwall/server
EnvironmentFile=/opt/codexwall/server/.env
ExecStart=/usr/bin/node /opt/codexwall/server/src/server.mjs
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/codexwall/server/data

[Install]
WantedBy=multi-user.target
